# ips.service
#!/bin/sh

[Unit]
Description=8006 Assignment 3

[Service]
ExecStart=/ips.sh /var/log/secure ssh 1 3

[Install]
WantedBy=multi-user.target
